# .config/systemd/user/clawdbot-gateway.service
[Unit]
Description=Clawdbot Gateway (v2026.1.25)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/npx clawdbot gateway --port 18789
WorkingDirectory=/home/liam
Restart=always
RestartSec=5
KillMode=control-group
Environment=HOME=/home/liam
Environment="PATH=/home/liam/.local/bin:/home/liam/.npm-global/bin:/home/liam/go-workspace/bin:/home/liam/.cargo/bin:/home/liam/bin:/home/liam/.nvm/current/bin:/home/liam/.fnm/current/bin:/home/liam/.volta/bin:/home/liam/.asdf/shims:/home/liam/.local/share/pnpm:/home/liam/.bun/bin:/usr/local/bin:/usr/bin:/bin"
Environment=CLAWDBOT_GATEWAY_PORT=18789
# Secrets loaded from secure env file (chmod 600)
EnvironmentFile=/home/liam/.clawdbot/credentials/gateway.env
Environment="CLAWDBOT_SYSTEMD_UNIT=clawdbot-gateway.service"
Environment=CLAWDBOT_SERVICE_MARKER=clawdbot
Environment=CLAWDBOT_SERVICE_KIND=gateway
Environment=CLAWDBOT_SERVICE_VERSION=2026.1.25

[Install]
WantedBy=default.target
