# ops/systemd/user/openclaw-gateway/openclaw-gateway.service.d/30-skills-sync.conf
[Service]
ExecStartPre=/usr/bin/bash -lc 'rsync -a --delete /home/dvk/openclaw-skills/ /home/dvk/clawd/skills/'
