# openclaw.service
[Unit]
Description=OpenClaw AI Gateway
After=network.target

[Service]
Type=simple
User=vova
WorkingDirectory=/home/vova/OpenPro
Environment=GEMINI_API_KEY=AIzaSyA2rEvczNGcf5FhScyVpOma57pqQghPEkg
Environment=NODE_ENV=production
Environment=OPENCLAW_STATE_DIR=/home/vova/OpenPro/workspace
Environment=PATH=/home/vova/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/node scripts/run-node.mjs gateway --force
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
