# runbook_memory/config/systemd/runbook-librarian.service
[Unit]
Description=OpenClaw runbook memory maintenance librarian
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/openclaw-upstream
ExecStart=%h/openclaw-upstream/runbook_memory/scripts/runbook-librarian.sh
