# scripts/systemd/chrome-cleanup.service
[Unit]
Description=OpenClaw Chrome Browser Process Cleanup
After=network.target
Documentation=https://github.com/openclaw/openclaw/blob/main/scripts/chrome-cleanup.sh

[Service]
Type=simple
# FIX: Updated path to match actual installation location
ExecStart=%h/.openclaw/workspace/scripts/chrome-cleanup.sh daemon
Restart=always
RestartSec=10
# FIX: Add graceful shutdown timeout
TimeoutStopSec=30
KillSignal=SIGTERM

[Install]
WantedBy=default.target
