# apps/pfos/orgo/systemd/pf-worker@.service
[Unit]
Description=Platinum Fang Worker %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/openclaw/apps/pfos
EnvironmentFile=%h/.config/pf-worker.env
Environment=PF_WORKER_ID=pf-worker-%i
ExecStart=/usr/bin/env npm run orchestrator:worker
Restart=always
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
