# scripts/email-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.email-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/pnpm</string>
        <string>exec</string>
        <string>tsx</string>
        <string>/Users/vero/openclaw/scripts/email-sync.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/vero/openclaw</string>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>/Users/vero/.openclaw/logs/email-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vero/.openclaw/logs/email-sync.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
