# scripts/morning-presync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.morning-presync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/vero/openclaw/scripts/morning-presync.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/vero/openclaw</string>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>50</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>50</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>50</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>50</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>50</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/vero/.openclaw/logs/morning-presync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vero/.openclaw/logs/morning-presync.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
