# sentinel/sentinel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clawd.sentinel</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/anaconda3/bin/python3</string>
        <string>/Users/sulaxd/clawd/sentinel/sentinel.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/sulaxd/clawd/sentinel</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/sulaxd/.local/bin:/opt/anaconda3/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/sulaxd/clawd/sentinel</string>
        <key>HOME</key>
        <string>/Users/sulaxd</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/sulaxd/clawd/sentinel/logs/sentinel-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/sulaxd/clawd/sentinel/logs/sentinel-stderr.log</string>

    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
