# examples/sofagenius-flyio/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:openclaw-gateway]
command=node /app/openclaw.mjs gateway --allow-unconfigured --port 3000 --bind lan
directory=/app
user=node
environment=NODE_ENV="production",HOME="/home/node"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=5

[program:sofagenius-backend]
command=python3 -m uvicorn backend.main:app --host 127.0.0.1 --port 8000 --app-dir /opt/sofagenius
directory=/opt/sofagenius
user=node
environment=HOME="/home/node"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=5
startsecs=5
